# docker/helix-local/supervisord.conf
[supervisord]
nodaemon=true

[program:p4d]
command=/usr/sbin/p4dctl start -a
